# deploy/99-wlan1-monitor.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlan1", \
  RUN+="/usr/local/bin/set-monitor-mode.sh"
