# deploy/kismet.service
[Unit]
Description=Kismet passive sensor daemon
After=network.target gpsd.service
Wants=gpsd.service

[Service]
User=kismet
Group=kismet
Environment=KISMET_CAPTURE_INTERFACE=wlan1
Environment=HOME=/var/lib/kismet
WorkingDirectory=/var/lib/kismet
ExecStart=/usr/bin/kismet --no-ncurses-wrapper -c wlan1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
