# deploy/passive-vigilance.service
[Unit]
Description=Passive Vigilance sensor orchestrator
After=network.target gpsd.service kismet.service
Wants=gpsd.service kismet.service
# Restart-loop backstop: if the service restarts StartLimitBurst times within
# StartLimitIntervalSec, systemd stops respawning it and leaves it failed instead
# of flailing forever. This is the outer guard for a systemd-watchdog SIGABRT loop
# (soak #3: a wedged alert backend starved the WATCHDOG=1 heartbeat, so systemd
# killed and respawned the node every ~2 min indefinitely; the in-process
# crash-guard only counts os._exit self-restarts, not systemd kills). A genuine
# one-off stall still self-heals — one restart, well under the burst — so only a
# non-recovering loop trips this and parks the unit for a human to look at.
StartLimitIntervalSec=900
StartLimitBurst=5

[Service]
User=${PI_USER}
Group=${PI_USER}
WorkingDirectory=/home/${PI_USER}/Passive-Vigilance
EnvironmentFile=/home/${PI_USER}/Passive-Vigilance/.env
# Enable Low Energy on the BLE controller before capture starts. After a reboot
# or USB re-plug the controller comes up BR/EDR-only (bluetoothd is intentionally
# not running, so nothing else turns LE on), which silently starves the raw-HCI
# advertisement scanner — it binds the socket and logs "started" but receives no
# LE adverts, so no BLE fingerprints are built. Retries to tolerate late USB
# enumeration and auto-detects the lowest-numbered controller (it may re-enumerate
# to hci1 after a reset). Non-fatal (leading '-') so a missing/absent controller
# never blocks orchestrator startup. Runs under the unit's CAP_NET_ADMIN, so
# btmgmt needs no sudo. Skipped entirely unless BLE_SCANNER_ENABLED=true.
ExecStartPre=-/bin/sh -c 'test "$BLE_SCANNER_ENABLED" = "true" || exit 0; for i in 1 2 3 4 5 6 7 8; do idx=$(ls /sys/class/bluetooth/ 2>/dev/null | sed -n "s/^hci\\([0-9][0-9]*\\)$/\\1/p" | sort -n | head -1); if [ -n "$idx" ]; then /usr/bin/btmgmt --index "$idx" power on >/dev/null 2>&1; /usr/bin/btmgmt --index "$idx" le on >/dev/null 2>&1 && exit 0; fi; sleep 2; done'
ExecStart=/opt/passive-vigilance/venv/bin/python3 /home/${PI_USER}/Passive-Vigilance/main.py
# Passive BLE capture (BLE_SCANNER_ENABLED=true) opens a raw HCI socket, which
# needs CAP_NET_RAW + CAP_NET_ADMIN. Granted as ambient capabilities so they apply
# to this non-root service process only — scoped here, not via setcap on the shared
# python binary. Harmless when BLE is disabled (the socket is never opened).
# Do NOT add CapabilityBoundingSet to restrict the set: it strips the SETUID/SETGID
# caps that `sudo` needs, which breaks the SDR coordinator's `sudo systemctl
# stop/start readsb` handoff (readsb then never releases the dongle, DroneRF can't
# claim it -> exitcode=0 crash-loop -> auto-disable). The default bounding set
# already contains the net caps, so AmbientCapabilities alone is sufficient.
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
# Type=notify: the orchestrator sends READY=1 once startup completes, then a
# WATCHDOG=1 heartbeat only while capture is genuinely progressing. If capture
# freezes, the heartbeats stop and systemd restarts the service after WatchdogSec
# — the outer net that survives even a total in-process wedge.
Type=notify
NotifyAccess=main
WatchdogSec=120
# READY=1 is sent only after startup() returns, and startup waits up to
# GPS_STARTUP_TIMEOUT_SECONDS (default 120) for a first fix. The default notify
# start timeout is 90s, so a slow-GPS cold boot would be killed before READY=1
# and crash-loop. Allow startup to exceed the GPS wait with margin.
TimeoutStartSec=180
Restart=always
RestartSec=5
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
