# deploy/passive-vigilance.service
[Unit]
Description=Passive Vigilance sensor orchestrator
After=network.target gpsd.service kismet.service
Wants=gpsd.service kismet.service

[Service]
User=${PI_USER}
Group=${PI_USER}
WorkingDirectory=/home/${PI_USER}/Passive-Vigilance
EnvironmentFile=/home/${PI_USER}/Passive-Vigilance/.env
ExecStart=/usr/bin/python3 /home/${PI_USER}/Passive-Vigilance/main.py
Restart=on-failure
RestartSec=10
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
